# alise.service.in
[Unit]
Description=alise
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=@EXE_PATH@
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
