# package/src/avm-alise/deb-scripts/postinst
cp /usr/local/share/avm-alise/avm-alise.service /lib/systemd/system/avm-alise.service
chmod 644 /lib/systemd/system/avm-alise.service
systemctl daemon-reload
systemctl enable avm-alise.service
systemctl start avm-alise.service
systemctl status avm-alise.service --no-pager || true

if pgrep -x avm-alise >/dev/null
then
    echo "$service_name started"
else
    echo "$service_name stopped"
fi
