# package/src/avm-alise/deb-scripts/prerm
systemctl disable avm-alise.service 2>/dev/null
systemctl stop avm-alise.service 2>/dev/null
