# package/src/avm-alise/deb-scripts/prerm
systemctl disable alise.service 2>/dev/null
systemctl stop alise.service 2>/dev/null