# package/src/avm-alise-stm/service/avm-alise.service
[Unit]
Description=avm-alise
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/local/bin/avm-alise
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
