# deb-scripts/postinst.in
deb_install_prefix=@CPACK_PACKAGING_INSTALL_PREFIX@
service_name=@SERVICE_NAME@
cp $deb_install_prefix/share/$service_name /lib/systemd/system/$service_name
chmod 644 /lib/systemd/system/$service_name
systemctl daemon-reload

systemctl enable $service_name
systemctl start $service_name
systemctl status $service_name --no-pager || true

if pgrep -x alise >/dev/null
then
    echo "$service_name started"
else
    echo "$service_name stopped"
fi
